# packages/network/pure-ftpd/system.d/ftpd.service
[Unit]
Description=pure-ftpd
After=network.target

[Service]
ExecStart=/usr/sbin/pure-ftpd
TimeoutStopSec=1
RestartSec=2
Restart=always
StartLimitInterval=4

[Install]
WantedBy=multi-user.target
