# packages/network/acestream/system.d/acestream.service
[Unit]
Description=AceStream Engine
Requires=network-online.target
After=network-online.target
Before=vdr.service kodi.service

[Service]
ExecStart=/opt/start-engine --client-console --live-cache-type memory --live-buffer 60 --live-cache-size 268435456 --upload-limit 1 --cache-limit 10
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=5
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
