# rc.local
#!/bin/bash

#e2fsck -fpDv -E bmap2extent /dev/sdb1
#e2fsck -fpDv -E bmap2extent /dev/sdc1
#e2fsck -fpDv -E bmap2extent /dev/sdd1

mount /dev/sdb1 /home/user
mount /dev/sdc1 /home/user/-f2fs
mount /dev/sdd1 /home/user/.aa

/bin/smbd
