# packages/thoradia/addons/aria2=/source/system.d/service.aria2.service
[Unit]
Description=aria2
Documentation=https://aria2.github.io/manual/en/html/index.html
After=multi-user.target network-online.target
Requires=multi-user.target network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.aria2/bin/aria2.start
Restart=on-failure

[Install]
WantedBy=kodi.target
