# packages/network/wireguard/system.d/wireguard-config.service
[Unit]
Description=WireGuard Config Service
After=local-fs.target
ConditionPathExists=/storage/.cache/wireguard/privatekey
ConditionPathExists=/storage/.cache/wireguard/publickey
ConditionPathExists=/storage/.cache/wireguard/preshared

[Service]
Type=oneshot
ExecStart=/usr/bin/wg-config
RemainAfterExit=yes
