# packages/network/xupnpd/system.d/xupnpd.service
[Unit]
Description=eXtensible UPnP agent
Requires=network-online.target
After=network-online.target
Wants=udpxy.service
Before=kodi.service
ConditionPathExists=/storage/.cache/services/xupnpd.conf

[Service]
ExecStartPre=-/usr/bin/xupnpd-config
ExecStart=/bin/sh -c '/usr/bin/xupnpd > /var/log/xupnpd.log'
PIDFile=/var/run/xupnpd.pid
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=3
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
