# packages/-my/shadowsocks-libev/i2pd.service
[Unit]
Description=Anonymizing overlay network for TCP
After=network.target

[Service]
Type=notify
NotifyAccess=all
PIDFile=/var/run/tor/tor.pid
PermissionsStartOnly=yes
ExecStart=/usr/bin/tor -f /storage/.config/tor/torrc --RunAsDaemon 0
ExecReload=/bin/kill -HUP ${MAINPID}
KillSignal=SIGINT
TimeoutSec=45
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
