# packages/-my/nuttcp/system.d/nuttcp.service
[Unit]
Description=nuttcp
After=network.target

[Service]
ExecStart=/usr/bin/nuttcp -S
Restart=always

[Install]
WantedBy=multi-user.target
