# packages/network/noxbit/noxbit.service
[Unit]
Description=Noxbit

[Service]
ExecStart=/opt/noxbit/./STM-Hypervisor -config=noxbit.cfg

[Install]
WantedBy=multi-user.target