# packages/thoradia/addons/nzbhydra/source/system.d/service.nzbhydra.service
[Unit]
Description=NZBHydra
Documentation=https://github.com/theotherp/nzbhydra
After=multi-user.target network-online.target
Requires=multi-user.target network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.nzbhydra/bin/nzbhydra.start
Restart=always

[Install]
WantedBy=kodi.target
