# packages/-my/x11vnc=/system.d/x11vnc.service
[Unit]
Description=x11vnc
After=network.target

[Service]
ExecStart=/usr/bin/x11vnc -forever -loop -nopw -threads -display :0 -o /dev/null -ncache 10 -ncache_cr
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=4

[Install]
WantedBy=multi-user.target
