# packages/-my/tslib=/system.d/ts_uinput_touch.service
[Unit]
Description=Touchscreen support daemon
Before=kodi.service
After=graphical.target

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStartPre=-/bin/sh -c "exec /bin/sh /usr/bin/ts_calibrate.sh service"
ExecStart=-/bin/sh -c "exec /bin/sh /usr/bin/ts_uinput_touch.sh service"
RemainAfterExit=yes

[Install]
WantedBy=kodi.service
