# packages/thoradia/addons/deluge/source/system.d/service.deluge.service
[Unit]
Description=Deluge service
Documentation=http://deluge-torrent.org/
After=multi-user.target network-online.target
Requires=multi-user.target network-online.target

[Service]
Type=oneshot
ExecStart=/bin/sh /storage/.kodi/addons/service.deluge/bin/deluge.start
ExecStop=/bin/sh /storage/.kodi/addons/service.deluge/bin/deluge.stop
RemainAfterExit=yes

[Install]
WantedBy=kodi.target
