# packages/-my/ps3remote/system.d/ps3remote.service
[Unit]
Description=Userspace PS3 Remote Driver
After=network.target

[Service]
ExecStart=/usr/sbin/ps3remote
Restart=always

[Install]
WantedBy=multi-user.target
