# packages/-my/ps3remote/system.d/ps3remote.service
[Unit]
Description=Userspace PS3 Remote Driver
After=network.target

[Service]
ExecStart=/usr/bin/ps3remote
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
StartLimitInterval=0

[Install]
WantedBy=kodi.target
