# packages/thoradia/addons/jackett/source/system.d/service.jackett.service
[Unit]
Description=Jackett
Documentation=https://github.com/Jackett/Jackett
After=multi-user.target network-online.target
Requires=multi-user.target network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.jackett/bin/jackett.start
Restart=always

[Install]
WantedBy=kodi.target
