# packages/network/iperf/system.d/iperf3.service
[Unit]
Description=iperf3
After=network.target

[Service]
ExecStart=/usr/bin/iperf3 -s
TimeoutStopSec=1
RestartSec=2
Restart=always
StartLimitInterval=4

[Install]
WantedBy=multi-user.target
