# packages/addons/service/snapclient/source/system.d/service.snapclient.service
[Unit]
Description=Snapclient
After=kodi.service
Wants=kodi.service

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.snapclient/bin/snapclient.start
Restart=always

[Install]
WantedBy=kodi.service
