# packages/addons/service/snapserver/source/system.d/service.snapserver.service
[Unit]
Description=Snapserver
After=multi-user.target
Before=kodi.service
Wants=multi-user.target kodi.service

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.snapserver/bin/snapserver.start
Restart=always

[Install]
WantedBy=multi-user.target kodi.service
