# packages/-my/x11vnc/system.d/x11vnc.service
[Unit]
Description=x11vnc
After=graphical.target
ConditionPathExists=/storage/.cache/services/x11vnc.conf

[Service]
ExecStart=/usr/bin/x11vnc.start
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
