# packages/-my/x11vnc/system.d/x11vnc.service
[Unit]
Description=x11vnc
After=network.target

[Service]
ExecStart=/usr/bin/x11vnc -forever -loop -threads -o /dev/null  -display :0
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=4

[Install]
WantedBy=multi-user.target
