# packages/thoradia/addons/sabnzbd/source/system.d/service.sabnzbd.service
[Unit]
Description=SABnzbd
Documentation=https://sabnzbd.org/
After=multi-user.target network-online.target
Requires=multi-user.target network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.sabnzbd/bin/sabnzbd.start
Restart=always

[Install]
WantedBy=kodi.target
