# packages/escalade/network/sonarr/system.d/sonarr.service
[Unit]
Description=Sonarr PVR for newsgroup and bittorrent users
Requires=network-online.service
After=network-online.service
ConditionPathExists=/storage/.kodi/addons/tools.mono/bin/mono

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=0
ExecStart=/storage/.kodi/addons/tools.mono/bin/mono /opt/sonarr/NzbDrone.exe -nobrowser

[Install]
WantedBy=multi-user.target
