# projects/Generic/squid_autostart25/squid.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--   Manifest-file for Squid, put this file in
       /var/svc/manifest/network/squid.xml
       and run #svccfg import /var/svc/manifest/network/squid.xml
       Fixed by Yuri Voinov (C) 2007,2015
-->
<service_bundle type='manifest' name='Squid'>
<service
        name='network/squid'
        type='service'
        version='1'>

        <create_default_instance enabled='false' />

        <single_instance />

        <dependency name='fs-local'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri
                        value='svc:/system/filesystem/local' />
        </dependency>

        <dependency name='net-loopback'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/network/loopback' />
        </dependency>

        <dependency name='net-physical'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/network/physical' />
        </dependency>

        <dependency name='utmp'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/system/utmp' />
        </dependency>

        <dependency name='squid_config_data'
                grouping='require_all'
                restart_on='refresh'
                type='path'>
                <service_fmri value='file://localhost/usr/local/squid/etc/squid.conf' />
        </dependency>

        <exec_method
                type='method'
                name='start'
                exec='/lib/svc/method/init.squid %m'
                timeout_seconds='900'/>

        <exec_method
                type='method'
                name='stop'
                exec='/lib/svc/method/init.squid %m'
                timeout_seconds='900' />

        <exec_method
                type='method'
                name='refresh'
                exec='/lib/svc/method/init.squid %m'
                timeout_seconds='180' />

        <exec_method
                type='method'
                name='restart'
                exec='/lib/svc/method/init.squid %m'
                timeout_seconds='180' />

        <property_group name='general' type='framework'>
                <!-- to start stop squid -->
                <propval name='action_authorization' type='astring'
                        value='solaris.smf.manage' />
        </property_group>

        <stability value='Unstable' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>
                        Squid proxy server
                        </loctext>
                </common_name>
                <documentation>
                        <manpage title='squid' section='8' manpath='/usr/local/squid/share/man/man8' />
                </documentation>
        </template>

</service>
</service_bundle>