# packages/-my/unbound/system.d/unbound.service
[Unit]
Description=unbound
After=network.target

[Service]
ExecStart=/usr/sbin/unbound -d
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=4

[Install]
WantedBy=multi-user.target
