# packages/thoradia/addons/sonarr/source/system.d/service.sonarr.service
[Unit]
Description=Sonarr
Documentation=https://sonarr.tv/
After=multi-user.target network-online.target
Requires=multi-user.target network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.sonarr/bin/sonarr.start
Restart=always

[Install]
WantedBy=kodi.target
