# xxx/-tvservice/tvheadend/system.d/tvheadend.service
[Unit]
Description=Tvheadend
After=network-online.target
Wants=network-online.target
ConditionPathExists=/storage/.cache/services/tvheadend.conf

[Service]
EnvironmentFile=-/run/alexelec/tvh.conf
ExecStartPre=-/usr/bin/tvh-config
ExecStartPre=-/usr/bin/camd start Tvheadend
ExecStart=/usr/bin/tvheadend $TVH_ARGS
ExecStop=/bin/kill -TERM $MAINPID
ExecStopPost=/usr/bin/camd stop Tvheadend
TimeoutStopSec=5
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
