# packages/thoradia/addons/radarr/source/system.d/service.radarr.service
[Unit]
Description=Radarr
Documentation=https://github.com/Radarr/Radarr
After=multi-user.target network-online.target
Requires=multi-user.target network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.radarr/bin/radarr.start
Restart=always

[Install]
WantedBy=kodi.target
